# dots/.config/systemd/user/qs-temps.service
[Unit]
Description=Write CPU/GPU temperatures for Quickshell

[Service]
Type=oneshot
ExecStart=%h/.local/bin/qs-temps
